# launchd-backup/com.investment-bot.monthly-deposit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- 월별 입금 cron LaunchAgent -->
<!-- 매월 1일 KST 09:00 (UTC 00:00)에 실행 -->
<plist version="1.0">
<dict>

    <!-- 서비스 식별자 -->
    <key>Label</key>
    <string>com.investment-bot.monthly-deposit</string>

    <!-- 실행 명령 -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jarvis/Projects/investment-bot/scripts/monthly_deposit_cron.py</string>
    </array>

    <!-- 작업 디렉토리 -->
    <key>WorkingDirectory</key>
    <string>/Users/jarvis/Projects/investment-bot</string>

    <!-- 환경변수 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- 실행 스케줄: 매월 1일 UTC 00:00 = KST 09:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- 상시 실행 불필요 (스케줄 기반) -->
    <key>RunAtLoad</key>
    <false/>

    <!-- 표준 출력/에러 → 같은 로그 파일 -->
    <key>StandardOutPath</key>
    <string>/Users/jarvis/Projects/investment-bot/logs/launchd_monthly.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jarvis/Projects/investment-bot/logs/launchd_monthly.log</string>

</dict>
</plist>
