# launchd/com.investment-bot.prices-krx.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- KRX + 미국장 가격 갱신 LaunchAgent -->
<!-- 10분(600초)마다 실행, 장중 여부는 refresh_prices.py가 스스로 판단하여 조기 종료 -->
<plist version="1.0">
<dict>

    <!-- 서비스 식별자 -->
    <key>Label</key>
    <string>com.investment-bot.prices-krx</string>

    <!-- 실행 명령 -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jarvis/Projects/investment-bot/scripts/refresh_prices.py</string>
    </array>

    <!-- 작업 디렉토리 -->
    <key>WorkingDirectory</key>
    <string>/Users/jarvis/Projects/investment-bot</string>

    <!-- 환경변수 (install.sh에서 실제 값으로 치환됨) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>KIWOOM_APPKEY</key>
        <string>__KIWOOM_APPKEY__</string>
        <key>KIWOOM_SECRETKEY</key>
        <string>__KIWOOM_SECRETKEY__</string>
        <key>BRAVE_API_KEY</key>
        <string>__BRAVE_API_KEY__</string>
        <!-- PATH: python3, 시스템 바이너리 접근용 -->
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- 실행 스케줄: 10분(600초)마다 반복 -->
    <key>StartInterval</key>
    <integer>600</integer>

    <!-- 표준 출력/에러 → 같은 로그 파일 -->
    <key>StandardOutPath</key>
    <string>/Users/jarvis/Projects/investment-bot/logs/launchd_prices_krx.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jarvis/Projects/investment-bot/logs/launchd_prices_krx.log</string>

    <!-- 상시 실행 불필요 (스케줄 기반) -->
    <key>RunAtLoad</key>
    <false/>

</dict>
</plist>
