# launchd/com.investment-bot.prices-us.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- 미국장 가격 갱신 LaunchAgent -->
<!-- 평일(월~금) KST 22:30~05:30, 매시간 실행 -->
<!-- KST → UTC: -9h -->
<!-- KST 22:30 = UTC 13:30 / KST 05:30(다음날) = UTC 20:30 -->
<!-- 미국장 KST 기준 월22:30~금05:30 → UTC 월13:30~금(+1일)20:30 -->
<!-- install.sh가 __PLACEHOLDER__ 값을 실제 환경변수로 치환함 -->
<plist version="1.0">
<dict>

    <!-- 서비스 식별자 -->
    <key>Label</key>
    <string>com.investment-bot.prices-us</string>

    <!-- 실행 명령 -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jarvis/Projects/investment-bot/scripts/refresh_prices.py</string>
    </array>

    <!-- 작업 디렉토리 -->
    <key>WorkingDirectory</key>
    <string>/Users/jarvis/Projects/investment-bot</string>

    <!-- 환경변수 (install.sh에서 실제 값으로 치환됨) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>KIWOOM_APPKEY</key>
        <string>__KIWOOM_APPKEY__</string>
        <key>KIWOOM_SECRETKEY</key>
        <string>__KIWOOM_SECRETKEY__</string>
        <key>BRAVE_API_KEY</key>
        <string>__BRAVE_API_KEY__</string>
        <!-- PATH: python3, 시스템 바이너리 접근용 -->
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- 실행 스케줄: 미국장 UTC 13:30~20:30 (KST 22:30~05:30) -->
    <!-- 월~금 UTC 13:30~19:30 + 화~토 UTC 20:30 (KST 다음날 05:30) -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- 월요일 UTC 13:30 = KST 22:30 -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 월요일 UTC 14:30 = KST 23:30 -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 월요일 UTC 15:30 = KST 24:30 (화 00:30) -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 월요일 UTC 16:30 = KST 화 01:30 -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 월요일 UTC 17:30 = KST 화 02:30 -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 월요일 UTC 18:30 = KST 화 03:30 -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 월요일 UTC 19:30 = KST 화 04:30 -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>

        <!-- 화요일 UTC 13:30 = KST 22:30 -->
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 화요일 UTC 14:30 = KST 23:30 -->
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 화요일 UTC 15:30 = KST 수 00:30 -->
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 화요일 UTC 16:30 = KST 수 01:30 -->
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 화요일 UTC 17:30 = KST 수 02:30 -->
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 화요일 UTC 18:30 = KST 수 03:30 -->
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 화요일 UTC 19:30 = KST 수 04:30 -->
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>

        <!-- 수요일 UTC 13:30 = KST 22:30 -->
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 수요일 UTC 14:30 = KST 23:30 -->
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 수요일 UTC 15:30 = KST 목 00:30 -->
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 수요일 UTC 16:30 = KST 목 01:30 -->
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 수요일 UTC 17:30 = KST 목 02:30 -->
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 수요일 UTC 18:30 = KST 목 03:30 -->
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 수요일 UTC 19:30 = KST 목 04:30 -->
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>

        <!-- 목요일 UTC 13:30 = KST 22:30 -->
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 목요일 UTC 14:30 = KST 23:30 -->
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 목요일 UTC 15:30 = KST 금 00:30 -->
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 목요일 UTC 16:30 = KST 금 01:30 -->
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 목요일 UTC 17:30 = KST 금 02:30 -->
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 목요일 UTC 18:30 = KST 금 03:30 -->
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 목요일 UTC 19:30 = KST 금 04:30 -->
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>

        <!-- 금요일 UTC 13:30 = KST 22:30 -->
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 금요일 UTC 14:30 = KST 23:30 -->
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 금요일 UTC 15:30 = KST 토 00:30 -->
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 금요일 UTC 16:30 = KST 토 01:30 -->
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 금요일 UTC 17:30 = KST 토 02:30 -->
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 금요일 UTC 18:30 = KST 토 03:30 -->
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 금요일 UTC 19:30 = KST 토 04:30 -->
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>

        <!-- 화요일 UTC 20:30 = KST 화 05:30 (월요일 미국장 마감) -->
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 수요일 UTC 20:30 = KST 수 05:30 (화요일 미국장 마감) -->
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 목요일 UTC 20:30 = KST 목 05:30 (수요일 미국장 마감) -->
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 금요일 UTC 20:30 = KST 금 05:30 (목요일 미국장 마감) -->
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 토요일 UTC 20:30 = KST 토 05:30 (금요일 미국장 마감) -->
        <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <!-- 표준 출력/에러 → 같은 로그 파일 -->
    <key>StandardOutPath</key>
    <string>/Users/jarvis/Projects/investment-bot/logs/launchd_prices_us.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jarvis/Projects/investment-bot/logs/launchd_prices_us.log</string>

    <!-- 상시 실행 불필요 (스케줄 기반) -->
    <key>RunAtLoad</key>
    <false/>

</dict>
</plist>
