# launchd/com.investment-bot.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.investment-bot.dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/jarvis/Projects/investment-bot/web/server.py</string>
  </array>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>WorkingDirectory</key>
  <string>/Users/jarvis/Projects/investment-bot</string>

  <key>StandardOutPath</key>
  <string>/Users/jarvis/Projects/investment-bot/logs/launchd_dashboard.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jarvis/Projects/investment-bot/logs/launchd_dashboard.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
