# launchd/com.investment-bot.news-morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- 장전 뉴스 갱신 LaunchAgent -->
<!-- 평일 KST 08:00 (UTC 전날 23:00) 실행 -->
<!-- KST 월~금 08:00 = UTC 일~목 23:00 (Weekday 0=일, 1=월, 2=화, 3=수, 4=목) -->
<!-- install.sh가 __PLACEHOLDER__ 값을 실제 환경변수로 치환함 -->
<plist version="1.0">
<dict>

    <!-- 서비스 식별자 -->
    <key>Label</key>
    <string>com.investment-bot.news-morning</string>

    <!-- 실행 명령 -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jarvis/Projects/investment-bot/scripts/refresh_news.py</string>
    </array>

    <!-- 작업 디렉토리 -->
    <key>WorkingDirectory</key>
    <string>/Users/jarvis/Projects/investment-bot</string>

    <!-- 환경변수 (install.sh에서 실제 값으로 치환됨) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>KIWOOM_APPKEY</key>
        <string>__KIWOOM_APPKEY__</string>
        <key>KIWOOM_SECRETKEY</key>
        <string>__KIWOOM_SECRETKEY__</string>
        <key>BRAVE_API_KEY</key>
        <string>__BRAVE_API_KEY__</string>
        <!-- PATH: python3, 시스템 바이너리 접근용 -->
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- 실행 스케줄: KST 월~금 08:00 = UTC 일~목 23:00 -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- 일요일 UTC 23:00 = KST 월요일 08:00 -->
        <dict>
            <key>Weekday</key>
            <integer>0</integer>
            <key>Hour</key>
            <integer>23</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 월요일 UTC 23:00 = KST 화요일 08:00 -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>23</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 화요일 UTC 23:00 = KST 수요일 08:00 -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>23</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 수요일 UTC 23:00 = KST 목요일 08:00 -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>23</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 목요일 UTC 23:00 = KST 금요일 08:00 -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>23</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <!-- 표준 출력/에러 → 같은 로그 파일 -->
    <key>StandardOutPath</key>
    <string>/Users/jarvis/Projects/investment-bot/logs/launchd_news_morning.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jarvis/Projects/investment-bot/logs/launchd_news_morning.log</string>

    <!-- 상시 실행 불필요 (스케줄 기반) -->
    <key>RunAtLoad</key>
    <false/>

</dict>
</plist>
